# docs/deploy/mediakit-server.service
[Unit]
Description=mediakit HTTP server
After=network.target redis.service
Requires=redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mediakit
EnvironmentFile=/home/ubuntu/mediakit/.env
ExecStart=/home/ubuntu/mediakit/.venv/bin/mediakit-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
